# systemd/backInTime-gpgKeys-Backup.service
[Unit]
Description=Backup das Chaves GPG

[Service]
ExecStart=/bin/bash ~/ShellScript/backInTime-gpgKeys-backup.sh
