# systemd/pipewire-reset.timer
[Unit]
Description=Restart PipeWire weekly

[Timer]
OnActiveSec=7d
Persistent=true

[Install]
WantedBy=timers.target
