# Viamar-PC/systemd/root-level/jkyon-paru-sync.service
[Unit]
Description=Paru database sync (pacman/paru)

[Service]
Type=oneshot
ExecStart=/usr/bin/paru -Syy --needed