# Viamar-PC/systemd/root-level/pacman-cache-update.service
[Unit]
Description=Atualiza o banco de dados do pacman/paru

[Service]
Type=oneshot
ExecStart=/usr/bin/paru -Sy --needed